# docker/cron/crontab
* * * * * cd /var/www/html && php artisan schedule:run >> /dev/stdout 1>&1
